# etc/profile.d/mybin.sh
export PATH=$PATH:/etc/my/bin

alias top=htop
alias df=pydf
alias grep='grep --colour=auto'
alias l='ls -la'
