# etc/profile.d/mybin.sh
echo "This system uses Romans My-Gentoo-Misc repsoitory"
echo "https://github.com/romaninsh/my-gentoo-misc"
export PATH=$PATH:/etc/my/bin
